# templates/systemd.service.j2
[Unit]
Description="{{ file_name }}"
After=network.target

[Service]
Type=forking
{% if environment_file is defined %}
EnvironmentFile={{ environment_file }}
{% endif %}
ExecStart={{ initscripts_path }}/{{ file_name }} start
ExecStop={{ initscripts_path }}/{{ file_name }} stop
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStartSec={{ timeout_start_sec }}
KillMode=process
Restart=on-failure
RestartSec=120
RestartPreventExitStatus=255
StartLimitBurst=0
{% if confluent_service.value.limit_nofile is defined %}
LimitNOFILE={{ limit_nofile }}
{% endif %}
{% if limit_nproc is defined %}
LimitNPROC={{ limit_nproc }}
{% endif %}

[Install]
WantedBy=multi-user.target
Alias={{ service_name }}.service
